# ansible/roles/rollback-guard/templates/travel-router-deadman.service.j2
[Unit]
Description=Deadman switch do travel-router (rollback se o heartbeat vencer)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/travel-router-deadman.sh
