# ansible/roles/access-point/templates/ap-watchdog.timer.j2
[Unit]
Description=Roda o watchdog do AP a cada minuto

[Timer]
OnBootSec=30
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target
