# ansible/roles/access-point/templates/ap-watchdog.service.j2
[Unit]
Description=Watchdog do AP privado (repara uap0/hostapd/IP/dnsmasq)
After=hostapd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ap-watchdog.sh
